# scripts/pimic-display.service
[Unit]
Description=PIMIC Dashboard Display Service
After=multi-user.target graphical.target
Wants=multi-user.target

[Service]
Type=simple
User=ochtii
Environment=DISPLAY=:0
Environment=HOME=/home/ochtii
WorkingDirectory=/home/ochtii/pimic
ExecStartPre=/bin/sleep 15
ExecStart=/home/ochtii/pimic/scripts/start_dashboard_display.sh
ExecStop=/usr/bin/pkill -f chromium
Restart=always
RestartSec=10
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target