# scripts/pimic-display.service
[Unit]
Description=PIMIC Dashboard Display Service
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=forking
User=ochtii
Environment=DISPLAY=:0
Environment=HOME=/home/ochtii
WorkingDirectory=/home/ochtii/pimic
ExecStartPre=/bin/sleep 10
ExecStart=/home/ochtii/pimic/scripts/start_dashboard_display.sh
ExecStop=/usr/bin/pkill -f chromium
Restart=always
RestartSec=5
TimeoutStartSec=30

[Install]
WantedBy=graphical-session.target